# launchd/com.zhuge.stock.teod.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zhuge.stock.teod</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yinlei/Desktop/量化/stock_screener/scripts/run_t_eod.sh</string>
    </array>

    <!-- 周一到周五 15:30（收盘后 30 分钟数据已结算）
         2026-06-01 引入：T 模块收盘汇总（配对 B/S + 算盈亏 + 写 t_summary JSON） -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>35</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>35</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>35</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>35</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>5</integer>
            <key>Hour</key><integer>15</integer>
            <key>Minute</key><integer>35</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
