# launchd/com.zhuge.stock.monthlyreview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zhuge.stock.monthlyreview</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yinlei/Desktop/量化/stock_screener/scripts/run_monthly_review_auto.sh</string>
    </array>

    <!-- 每月 1 号 17:00 跑（统计上月数据，无视星期几）
         2026-06-01 引入：自动月复盘，告别手动 bash run_monthly_review.sh -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key><integer>1</integer>
        <key>Hour</key><integer>17</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
