# launchd/com.zhuge.stock.summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zhuge.stock.summary</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/yinlei/Desktop/量化/stock_screener/scripts/run_summary.sh</string>
    </array>

    <!-- 每周五 19:10（改自 15:40，与 update_review 19:00 排程对齐：等 T+1 数据补全后再跑） -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>5</integer>
        <key>Hour</key><integer>19</integer>
        <key>Minute</key><integer>10</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
